# .config/systemd/user/graphical-login.target
[Unit]
Description=Target started with i3

[Install]
WantedBy=default.target
